# deploy/systemd/upwork-agent-browser-session.service
[Unit]
Description=Upwork Agent browser session
After=network.target
Wants=network.target

[Service]
Type=simple
User=upwork-agent
WorkingDirectory=/opt/upwork-agent/app
EnvironmentFile=-/opt/upwork-agent/app/.env
Environment=NODE_ENV=production
Environment=DISPLAY=:1
ExecStart=/usr/bin/env bash scripts/start-browser-session.sh
Restart=on-failure
RestartSec=10
NoNewPrivileges=true
PrivateTmp=true
ProtectHome=read-only
ReadWritePaths=/opt/upwork-agent/app /opt/upwork-agent/shared /home/upwork-agent/.vnc

[Install]
WantedBy=multi-user.target
